# bscheduler.service.in
[Unit]
Description=Bscheduler daemon
After=network.target

[Service]
EnvironmentFile=-/etc/default/@project_name@
ExecStart=@prefix@/@sbindir@/@bscheduler_exe@ $BSCHEDULER_ARGS
Restart=always
RestartSec=7

[Install]
WantedBy=multi-user.target
